# pkgs/50-raspion/raspion-1.1.0/debian/postinst
#!/bin/sh
# postinst script for raspion
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# source default pkg values
. /usr/share/raspion/etc.templates/default.conf

apply_etctemplate() {
    local CONF=$1
    local NAME=$(basename $CONF)
    [[ -f /usr/share/raspion/etc.templates/$NAME ]] || return 1
    cp /usr/share/raspion/etc.templates/$NAME $CONF
    #manual sed chaos - replace global
    sed -i "s/#IPv6HOST#/$IPV6HOST/g" $CONF
    sed -i "s/#IPv6NET#/$IPV6NET/g" $CONF
    sed -i "s/#Pv4HOST#/$IPV4HOST/g" $CONF
    sed -i "s/#IPv4NET#/$IPV4NET/g" $CONF
    sed -i "s/#DHCPv4START#/$DHCPv4START/g" $CONF
    sed -i "s/#DHCPv4END#/$DHCPv4END/g" $CONF
    sed -i "s/#SSID#/$SSID/g" $CONF
    return 0
}

case "$1" in
    configure)
    #add pi user to wireshark and www-data
    usermod -a -G wireshark pi
    usermod -a -G www-data pi
    
    #add iptable rules
    iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE
    ip6tables -t nat -A POSTROUTING -o eth0 -s $IPV6NET/64 -j MASQUERADE 
    iptables -A PREROUTING -t nat -p tcp --dport 80 -j REDIRECT --to-ports 81 -i eth0
    ip6tables -A PREROUTING -t nat -p tcp --dport 80 -j REDIRECT --to-ports 81 -i eth0 
    #save policys
    netfilter-persistent save

    #modify & place etc.templates
    apply_etctemplate "/etc/ntopng/ntopng.conf"
    apply_etctemplate "/etc/gtk-3.0/settings.ini"
    apply_etctemplate "/etc/default/shellinabox"
    apply_etctemplate "/etc/network/interfaces"
    apply_etctemplate "/etc/hostapd/hostapd.conf"
    apply_etctemplate "/etc/radvd.conf"

    #generate hostapd pw and append to conf
    WPAPW=$(pwgen --ambiguous 9) 
    echo "wpa_passphrase=$WPAPW" >> /etc/hostapd/hostapd.conf.raspion
    echo "The hostapd Password ist: $WPAPW"
    
    #change hostname to raspion
    sed -i "s/^.*127.0.1.1*/127.0.1.1   raspion/" /etc/hosts
    echo "raspion" > /etc/hostname
    
    #unmask hostapd and enable radvd
    systemctl unmask hostapd
    systemctl enable radvd

    #set public_html permissions
    chmod g+s /home/pi/public_html/caps
    chmod 777 /home/pi/public_html/caps
    chgrp www-data /home/pi/public_html/caps

    
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    echo "something really bad happend\n - take a break\n - save modified confs\n - and rerun configuration"
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
