# pkgs/10-raspion/raspion-1.1.0/debian/postinst
#!/bin/sh
##
# c't-Raspion, a Raspberry Pi based all in one sniffer
# for judging on IoT and smart home devices activity
# (c) 2019-2020 c't magazin, Germany, Hannover
# see: https://ct.de/-123456 for more information
# 
# file          : postinst
# description   : postinst for raspion
# version       : 0.1
#       
# changes         
# author		: Benny Stark - github.com/Diggen85
# date          : 2020010
# notes         : Initial - untested
#               : debhelper dh_installdeb (1)
# date          : 2020119
# notes         : do pkg-divert without config-package-dev
###

set -e
#set -x

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# source default pkg values
. /usr/share/raspion/etc.templates/default.conf


# Config-package-dev displace/divert does not behave as need
# divert while customizing the templates should do a better job
# the code is inspired from config-package-dev scripts - BST200119
apply_etctemplate() {
    local PKG_NAME="raspion"
    local CONF_FILE=$1
    local CONF_NAME=$(basename $1)
    local SED_SCRIPT="s/#IPV6HOST#/$IPV6HOST/g;\
        s/#IPV6NET#/$IPV6NET/g;\
        s/#PV4HOST#/$IPV4HOST/g;\
        s/#IPV4NET#/$IPV4NET/g;\
        s/#DHCPV4START#/$DHCPv4START/g;\
        s/#DHCPV4END#/$DHCPv4END/g;\
        s/#SSID#/$SSID/g;"

    #check if the template exists
    test -f /usr/share/raspion/etc.templates/$CONF_NAME || return 1
    #test if we divert already or must divert
    LC_ALL=C dpkg-divert --list "$PKG_NAME" | grep -xiqe "^diversion.*$CONF_FILE.*$PKG_NAME.*$" || \

    #test if there is already a conf and backup
    test -f $CONF_FILE && mv -f $CONF_FILE $CONF_FILE.$PKG_NAME-bak
    
    #run sed on the template and write  config file
    sed -e "$SED_SCRIPT" /usr/share/raspion/etc.templates/$CONF_NAME > $CONF_FILE

    return 0
}

case "$1" in
    configure)
        #add pi user to wireshark and www-data
        usermod -a -G wireshark pi
        usermod -a -G www-data pi
        
        #clean and set raspion iptable rules
        iptables -X
        ip6tables -X
        iptables -F
        ip6tables -F
        iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE
        ip6tables -t nat -A POSTROUTING -o eth0 -s $IPV6NET/64 -j MASQUERADE 
        iptables -A PREROUTING -t nat -p tcp --dport 80 -j REDIRECT --to-ports 81 -i eth0
        ip6tables -A PREROUTING -t nat -p tcp --dport 80 -j REDIRECT --to-ports 81 -i eth0 
        #save policys
        netfilter-persistent save

        #modify & place etc.templates
        apply_etctemplate "/etc/ntopng/ntopng.conf"
        apply_etctemplate "/etc/gtk-3.0/settings.ini"
        apply_etctemplate "/etc/default/shellinabox"
        apply_etctemplate "/etc/network/interfaces"
        apply_etctemplate "/etc/hostapd/hostapd.conf"
        apply_etctemplate "/etc/radvd.conf"
        apply_etctemplate "/etc/lighttpd/conf-available/20-extport.conf"
        
        # Enable lighthttpd config
        ln -sf /etc/lighttpd/conf-available/20-extport.conf /etc/lighttpd/conf-enabled/20-extport.conf
        systemctl restart lighthttpd

        #generate hostapd pw and append to conf
        WPAPW=$(pwgen --ambiguous 9) 
        echo "wpa_passphrase=$WPAPW" >> /etc/hostapd/hostapd.conf
        echo "The hostapd Password is: $WPAPW"
        unset $WPAPW

        #change hostname to raspion
        sed -i "s/^.*127\.0\.1\.1.*$/127.0.1.1   raspion/" /etc/hosts
        echo "raspion" > /etc/hostname
        
        #unmask hostapd and enable radvd
        systemctl unmask hostapd
        systemctl enable radvd

        #set permissions for pi's (wireshark) .config-dir, installed by our deb
        chown pi:pi /home/pi/.config
        chown -R pi:pi /home/pi/.config/wireshark

        #set public_html permissions
        chown -R pi:pi /home/pi/public_html
        chmod g+s /home/pi/public_html/caps
        chmod 777 /home/pi/public_html/caps
        chgrp www-data /home/pi/public_html/caps
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        echo "something bad happend\n - take a break\n - better check confs and rerun configuration"
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
